# conf/survey.conf
description "surveys.ecit.fsu.edu"
author "ecitlabs@gmail.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env NEW_RELIC_CONFIG_FILE=/usr/share/nginx/surveys.ecit.fsu.edu/newrelic.ini
chdir /usr/share/nginx/surveys.ecit.fsu.edu
exec /root/.venvs/survey/bin/newrelic-admin run-program /root/.venvs/survey/bin/gunicorn --bind localhost:5000 --workers 2 app:app --log-level=debug --log-file=/var/log/gunicorn/app.log --preload